# scripts/systemd/foritech-pqc-monitor.timer
[Unit]
Description=Every-minute health check for Foritech PQC

[Timer]
OnBootSec=30s
OnUnitActiveSec=60s
AccuracySec=10s
Unit=foritech-pqc-monitor.service

[Install]
WantedBy=timers.target
